# deploy/quadlet/hyperboard-api.container
[Unit]
Description=Hyperboard API
After=hyperboard-postgresql.service hyperboard-rustfs.service
Requires=hyperboard-postgresql.service hyperboard-rustfs.service

[Container]
Image=localhost/hyperboard-api:latest
ContainerName=hyperboard-api
Network=hyperboard.network
Environment=HYPERBOARD_API_POSTGRESQL_HOST=hyperboard-postgresql
Environment=HYPERBOARD_API_POSTGRESQL_USER=hyperboard
Environment=HYPERBOARD_API_POSTGRESQL_PASSWORD=hyperboard
Environment=HYPERBOARD_API_POSTGRESQL_DATABASE=hyperboard
Environment=HYPERBOARD_API_POSTGRESQL_SSL_MODE=disable
Environment=HYPERBOARD_API_ADMIN_PASSWORD=hyperboard
Environment=HYPERBOARD_API_STORAGE_ENDPOINT=http://hyperboard-rustfs:9000
Environment=HYPERBOARD_API_STORAGE_BUCKET=hyperboard
Environment=HYPERBOARD_API_STORAGE_ACCESS_KEY=rustfsadmin
Environment=HYPERBOARD_API_STORAGE_SECRET_KEY=rustfsadmin
Environment=HYPERBOARD_API_STORAGE_REGION=us-east-1
Environment=HYPERBOARD_API_STORAGE_USE_PATH_STYLE=true
Environment=HYPERBOARD_API_LOG_LEVEL=info
PublishPort=8081:8080
HealthCmd=curl -f http://localhost:8080/healthz || exit 1
HealthInterval=5s
HealthTimeout=3s
HealthRetries=5

[Service]
Restart=always

[Install]
WantedBy=multi-user.target default.target
