# deploy/quadlet/hyperboard-api.container
[Unit]
Description=Hyperboard API
After=hyperboard-postgresql.service hyperboard-rustfs.service
Requires=hyperboard-postgresql.service hyperboard-rustfs.service

[Container]
Image=ghcr.io/dharmab/hyperboard-api:main
ContainerName=hyperboard-api
Pod=hyperboard.pod
Environment=HYPERBOARD_API_PORT=8081
Environment=HYPERBOARD_API_POSTGRESQL_HOST=localhost
Environment=HYPERBOARD_API_POSTGRESQL_USER=${HYPERBOARD_DB_USER}
Environment=HYPERBOARD_API_POSTGRESQL_PASSWORD=${HYPERBOARD_DB_PASSWORD}
Environment=HYPERBOARD_API_POSTGRESQL_DATABASE=${HYPERBOARD_DB_NAME}
Environment=HYPERBOARD_API_POSTGRESQL_SSL_MODE=disable
Environment=HYPERBOARD_API_ADMIN_PASSWORD=${HYPERBOARD_ADMIN_PASSWORD}
Environment=HYPERBOARD_API_STORAGE_ENDPOINT=http://localhost:9000
Environment=HYPERBOARD_API_STORAGE_BUCKET=${HYPERBOARD_STORAGE_BUCKET}
Environment=HYPERBOARD_API_STORAGE_ACCESS_KEY=${HYPERBOARD_STORAGE_ACCESS_KEY}
Environment=HYPERBOARD_API_STORAGE_SECRET_KEY=${HYPERBOARD_STORAGE_SECRET_KEY}
Environment=HYPERBOARD_API_STORAGE_REGION=${HYPERBOARD_STORAGE_REGION}
Environment=HYPERBOARD_API_STORAGE_USE_PATH_STYLE=true
Environment=HYPERBOARD_API_LOG_LEVEL=info
HealthCmd=curl -f http://localhost:8081/healthz || exit 1
HealthInterval=5s
HealthTimeout=3s
HealthRetries=5

[Service]
EnvironmentFile=/etc/hyperboard/hyperboard.env
Restart=always
